# files/easypanel/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

; ---- Redis Main (Enketo persistent store) ----
[program:redis-main]
command=redis-server /etc/redis/redis-main.conf --dir /var/lib/redis/main
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ---- Redis Cache (Enketo XSLT cache) ----
[program:redis-cache]
command=redis-server /etc/redis/redis-cache.conf --dir /var/lib/redis/cache
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ---- PyXForm HTTP Server (XLSForm to XML converter) ----
[program:pyxform]
command=python3 /opt/pyxform/pyxform-server.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

; ---- Enketo (Web Forms Engine) ----
[program:enketo]
command=/bin/bash /srv/src/enketo/packages/enketo-express/start-enketo.sh
directory=/srv/src/enketo/packages/enketo-express
environment=ENKETO_SRC_DIR="/srv/src/enketo/packages/enketo-express",NODE_TLS_REJECT_UNAUTHORIZED="0"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
startsecs=5

; ---- ODK Central Backend (Node.js API) ----
[program:service]
command=/bin/bash /usr/odk/start-odk.sh
directory=/usr/odk
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=40
startsecs=10

; ---- Nginx (Frontend + Reverse Proxy) ----
[program:nginx]
command=/scripts/setup-odk.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=50
startsecs=3
